# launchd/com.local.memory.index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.memory.index</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/make</string>
    <string>-C</string>
    <string>/ABSOLUTE/PATH/TO/MemoryOS</string>
    <string>rebuild</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MEMORY_HOME</key>
    <string>/ABSOLUTE/PATH/TO/Memory</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/Memory/_system/logs/launchd-index.out.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/Memory/_system/logs/launchd-index.err.log</string>
</dict>
</plist>
