# launchd/com.memoryos.import-pending.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.memoryos.import-pending</string>
  <key>ProgramArguments</key>
  <array>
    <string>__MEMORYOS_PYTHON__</string>
    <string>-m</string>
    <string>memoryos.cli</string>
    <string>import-pending</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__MEMORYOS_REPOSITORY__</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__MEMORYOS_HOME__/_system/logs/import-pending.out.log</string>
  <key>StandardErrorPath</key>
  <string>__MEMORYOS_HOME__/_system/logs/import-pending.err.log</string>
</dict>
</plist>
